# distro/etc/systemd/system/boot.mount
[Unit]
Description=Boot Partition

[Mount]
What=/dev/vda1
Where=/boot
Type=vfat
Options=defaults,rw,sync

[Install]
WantedBy=multi-user.target
