# distro/etc/systemd/system/daemon-init.service
[Unit]
Description=daemon-init

[Service]
Type=oneshot
ExecStart=/etc/framework/daemons/init.sh

[Install]
WantedBy=multi-user.target